# external/overlay/etc/cron/crontabs/root
*/10 * * * * /etc/init.d/S48sntp start >/dev/null 2>&1
